# ansible_collections/yamisskey/servers/roles/db_backup/templates/db-backup.timer.j2
[Unit]
Description=Database Backup Timer

[Timer]
OnCalendar={{ db_backup_schedule }}
RandomizedDelaySec={{ db_backup_schedule_randomized_delay }}
Persistent=true

[Install]
WantedBy=timers.target
